# SOURCES/solr.init.in
#!/bin/bash
#-----------------------------------------------------
# solr-initd v0.5.1
#-----------------------------------------------------
# Script for running solr as a service under initd.
#
# Requires solr.conf and solr.start
# Usage: service solr {start|stop|restart|status}"
#
# NOTE: This script was build specifically for Fedora,
# but also works on RHEL and CentOS. Other distributions
# are untest and probably require you to alter the scripts
#-----------------------------------------------------
# Author: Jeff Busby <jeff@jeffbusby.ca>
# License: http://creativecommons.org/licenses/by/3.0/
#-----------------------------------------------------
### BEGIN INIT INFO
# Provides:          solr
# Required-Start:   
# Required-Stop:
# Default-Start:     3 5
# Default-Stop:      0 1 2 6
# Short-Description: init script for solr search
# Description:       init script for solr search
#                    This should be placed in /etc/init.d
### END INIT INFO

. /etc/rc.d/init.d/functions
. /etc/sysconfig/solr

NAME=Solr

if [ ! -f $SOLR_STOP_KEY ];then
  action "Generating Solr Stop Key: $SOLR_STOP_KEY"
  tr -dc "[:alpha:]" < /dev/urandom | head -c 8 > $SOLR_STOP_KEY
fi

START_COMMAND="cd ${SOLR_RUNDIR} && java $OPTIONS -DSTOP.KEY=`cat $SOLR_STOP_KEY` -jar start.jar"
STOP_COMMAND="cd ${SOLR_RUNDIR} && java $OPTIONS -DSTOP.KEY=`cat $SOLR_STOP_KEY` -jar start.jar --stop"

case $1 in
    start)
        action "Starting ${NAME}: " daemon --user "$SOLR_USER" --pidfile $PIDFILE $START_COMMAND &
        ;;
    stop)
        action "Stopping ${NAME}: " eval $STOP_COMMAND
        ;;
    restart)
        $0 stop
        $0 start
        ;;
    status)
        status solr
        ;;
    *)
        echo "Usage: $0 {start|stop|restart|status}"
        exit 3
        ;;
esac
